# files/usr/lib/systemd/system/waydroid-container.service.d/mirror-os.conf
[Service]
# Patch waydroid_base.prop with GBM/Mesa GPU props before every start so
# that hardware acceleration is not lost to a software fallback after reboots.
ExecStartPre=/usr/libexec/mirror-os/waydroid-prestart.sh
