# systemd/system/mirror-nix-setup.service
[Unit]
Description=Mirror OS — Home Manager first-boot setup
Documentation=https://mirror-os.com/docs/home-manager
# Must run after Nix is installed and a network connection is available,
# since home-manager switch will pull nixpkgs and flake inputs on first run.
Requires=install-nix.service
After=install-nix.service
After=network-online.target
Wants=network-online.target
# Only run in a graphical session context — the script needs a real user.
After=graphical.target

[Service]
Type=oneshot
# Run as root so the script can detect the logged-in user, scaffold files
# into /home, modify /etc/shells, and call chsh.
User=root
ExecStart=/usr/libexec/mirror-os/mirror-nix-setup.sh
# Keep the service state visible after it completes so that
# `systemctl status mirror-nix-setup` shows the last run result.
RemainAfterExit=yes
# Give the first-time Nix evaluation (nixpkgs + HM + nix-flatpak + nixvim
# downloads) enough time to complete on a slow connection.
TimeoutStartSec=1800
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mirror-nix-setup

[Install]
WantedBy=graphical.target
