# files/systemd/system/install-nix.service
[Unit]
Description=Install Nix Package Manager (first boot)
ConditionPathExists=!/var/lib/mirror-os/nix-installed
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/mirror-os/install-nix.sh
ExecStartPost=/bin/bash -c 'mkdir -p /var/lib/mirror-os && touch /var/lib/mirror-os/nix-installed'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target