# files/systemd/system/install-home-manager.service
[Unit]
Description=Mirror OS – first-boot Home Manager installation
Documentation=https://mirror-os.com
After=install-nix.service network-online.target
Requires=install-nix.service
Wants=network-online.target
ConditionPathExists=!/var/lib/mirror-os/home-manager-installed

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/mirror-os/install-home-manager.sh
ExecStartPost=/bin/bash -c 'mkdir -p /var/lib/mirror-os && touch /var/lib/mirror-os/home-manager-installed'
StandardOutput=journal
StandardError=journal
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target