# files/systemd/system/mirror-user-setup.service
[Unit]
Description=Mirror OS — per-user Home Manager first-login setup
Documentation=https://mirror-os.com/docs/home-manager
# Only run if Home Manager has not been scaffolded yet for this user.
# %h expands to the user's home directory.
ConditionPathExists=!%h/.config/home-manager/flake.nix
# Nix must be installed at system level before this runs.
# Network needed for nix-community/home-manager flake fetch.
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/mirror-os/mirror-user-setup.sh
RemainAfterExit=yes
# First run downloads nixpkgs, home-manager, nix-flatpak, nixvim, flatpaks.
# Allow up to 30 minutes on slow connections.
TimeoutStartSec=1800
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mirror-user-setup

[Install]
WantedBy=default.target