# files/etc/systemd/user/mirror-os-apps.path
[Unit]
Description=Mirror OS flatpak app change detector
After=mirror-os-init.service

[Path]
PathChanged=%h/.local/share/flatpak/app
MakeDirectory=true

[Install]
WantedBy=default.target
