# files/etc/systemd/user/mirror-os-sync.timer
[Unit]
Description=Mirror OS hourly state sync

[Timer]
OnBootSec=300s
OnCalendar=*:00:00
Persistent=true

[Install]
WantedBy=timers.target
