# files/etc/systemd/user/mirror-os-hm.service
[Unit]
Description=Mirror OS home-manager sync (triggered by config change)
After=mirror-os-init.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/mirror-os/mirror-sync
StandardOutput=journal
StandardError=journal
TimeoutStartSec=600
